# etc/systemd/makamujo.service
[Unit]
Description=makamujo service
After=graphical.target xorg10.service x11vnc-10.service
Wants=graphical.target xorg10.service x11vnc-10.service

# Service runs multiple backgrounded processes and the start script exits.
# Use oneshot + RemainAfterExit so systemd considers the service active
# while background processes run and ExecStop can stop them.
[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/makamujo
User=root
Environment=PATH=/root/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=DISPLAY=:10
ExecStart=/opt/makamujo/bin/start-with-xauth.sh
ExecStop=/opt/makamujo/bin/stop
StandardOutput=journal
StandardError=journal
# If you want to run it as a different user, change `User=` and set `XAUTHORITY` accordingly:
# User=yourusername
# Environment=DISPLAY=:0
# Environment=XAUTHORITY=/home/yourusername/.Xauthority

[Install]
WantedBy=graphical.target
