# etc/systemd/makamujo-screen.service
[Unit]
Description=makamujo screen service
After=makamujo.service
PartOf=makamujo.service

[Service]
Type=simple
WorkingDirectory=@PREFIX@
User=root
Environment=PATH=@PREFIX@/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=DISPLAY=:10
# Environment=CONSOLE_BASIC_AUTH_PASSWORD=...
ExecStart=/bin/sh -lc "cd '@PREFIX@' && NODE_ENV=production '@BUN_BIN@' start"
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5s
