# etc/systemd/xorg10.service
[Unit]
Description=Persistent Xorg server on :10
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Environment=DISPLAY=:10
ExecStart=/opt/makamujo/bin/xorg10
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
